# contrib/com.user.lola.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  lola LaunchAgent. Install into ~/Library/LaunchAgents (NOT
  /Library/LaunchDaemons — lola must run in your user/GUI session so its
  tmux-hosted agent sessions have GUI context and the keychain is reachable),
  then:

      launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.user.lola.plist

  Before loading, replace:
    - /usr/local/bin/lola  with the absolute path to your lola binary
    - /Users/YOU          with your real home directory (launchd does NOT
                          expand ~, so every path below must be absolute)
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.lola</string>

	<!-- Adjust the binary path to wherever you installed lola. -->
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/lola</string>
		<string>run</string>
	</array>

	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>

	<!-- launchd provides no login-shell environment: inject a PATH that can
	     find tmux/git/gh/claude, and HOME so ~/.lola resolves. -->
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/YOU</string>
	</dict>

	<key>WorkingDirectory</key>
	<string>/Users/YOU</string>

	<!-- launchd does not expand ~ — use the absolute path. -->
	<key>StandardOutPath</key>
	<string>/Users/YOU/.lola/daemon.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOU/.lola/daemon.log</string>
</dict>
</plist>
